# gv-grok-bridge.service
[Unit]
Description=GV-Grok Bridge
After=network-online.target systemd-user-sessions.service
Wants=network-online.target

[Service]
Type=simple
User=gvgrok
Group=gvgrok
WorkingDirectory=/home/gvgrok/project
Environment="HOME=/home/gvgrok"
Environment="USER=gvgrok"
Environment="XDG_RUNTIME_DIR=/run/user/1003"
Environment="GV_HEADLESS=true"
Environment="GV_DISPLAY_NUM=:99"
Environment="GV_AUTO_ACCEPT=true"
Environment="GV_LOG_LEVEL=info"
Environment="GV_AUTHORIZED_NUMBERS=+1XXXXXXXXXX"
Environment="PATH=/usr/local/bin:/usr/bin:/bin:/usr/local/sbin:/usr/sbin:/sbin"
ExecStartPre=-/usr/bin/bash -c "/usr/bin/pkill -9 -f 'chromium' 2>/dev/null || true"
ExecStart=/usr/bin/node /home/gvgrok/project/dist/main.js
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/gvgrok/project/bridge.log
StandardError=append:/home/gvgrok/project/bridge.log

[Install]
WantedBy=multi-user.target
